# device-jenkins-el9/src/main/resources/etc/systemd/system/jenkins.service.d/local.conf
[Unit]
After=jenkins-autodiscovery.service
Requisite=jenkins-autodiscovery.service

[Service]
EnvironmentFile = -/etc/sysconfig/jenkins
UMask=0002

